# puppet/modules/orthanc/files/orthanc.conf
description "Orthanc Server"
author "Fernando Serran"

start on runlevel [2345]
stop on shutdown

exec Orthanc /etc/orthanc/orthanc.json