# puppet/modules/orthanc/files/orthanc.conf
description "Orthanc Server"
author "Fernando Serrano"

start on runlevel [2345]
stop on shutdown
respawn

exec Orthanc /etc/orthanc/orthanc.json