# puppet/modules/orthanc/files/orthanc.conf
description "Orthanc Server"
author "Fernando Serran"

start on started postgresql
stop on shutdown

exec Orthanc /etc/orthanc/orthanc.json